# Config/etc/systemd/system/systemd-resolved.service.d/override.conf
[Unit]
After    = var.mount
Requires = var.mount
